# deploy/supervisor/chat-next.conf
[program:chat-next]
directory=/var/www/chat-frontend
command=node .next/standalone/server.js
user=www-data
autostart=true
autorestart=true
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/chat-next.log
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1"
