# deploy/supervisor/chat-reverb.conf
[program:chat-reverb]
directory=/var/www/chat-api
command=php artisan reverb:start --host=127.0.0.1 --port=8080
user=www-data
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/chat-reverb.log
