# deploy/supervisor/chat-horizon.conf
[program:chat-horizon]
directory=/var/www/chat-api
command=php artisan horizon
user=www-data
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/chat-horizon.log
