# deploy/supervisor/livekit.conf
[program:livekit]
directory=/var/lib/livekit
command=/usr/local/bin/livekit-server --config /etc/livekit/livekit.yaml
user=www-data
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/livekit.log
