# backend/deploy/systemd/laravel-scheduler.timer.example
[Unit]
Description=Run DDRBit Laravel scheduler every minute

[Timer]
OnBootSec=30
OnUnitActiveSec=60
Unit=laravel-scheduler.service

[Install]
WantedBy=timers.target
